# web08/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND

[program:php-server]
command=php -S 0.0.0.0:5000 /var/www/html/vuln132321321321.php
